# ug-paneld.service
[Unit]
Description=Ugreen NAS Front Panel Display
# Start early (right after basic system init), NOT after multi-user.target —
# the latter waits for every other service (ZFS import, VM autostart, slow
# disks, ...), which on a NAS can be minutes, leaving the panel dark until
# then. The daemon waits for the DRM connector itself (drm_probe_timeout), so
# an early start is safe.
After=basic.target

[Service]
Type=simple
# i2c-hid unbinding is handled inside ug-paneld itself (supports both the
# CUST0000:00 and MSFT8000:00 touchscreen ids and logs what it does).
# drivetemp provides SATA disk temperatures, i2c-dev the touchscreen bus.
ExecStartPre=-/sbin/modprobe drivetemp
ExecStartPre=-/sbin/modprobe i2c-dev
ExecStart=/usr/bin/ug-paneld
Restart=always
RestartSec=5
# Exit code 2 means the kernel exposed no connected DRM connector (internal
# panel missing/disabled). Restarting cannot fix that, so don't loop.
RestartPreventExitStatus=2

[Install]
WantedBy=multi-user.target
